# systemd/.config/systemd/user/rbackup.path
[Unit]
Description=Checks if paths that are currently being backed up have changed

[Path]
PathChanged=%h/org

[Install]
WantedBy=default.target
