# systemd/.config/systemd/user/rbackup.service
[Unit]
Description=Backs up files

[Service]
ExecStart=/usr/bin/rsync %h/./org -CERrltm --delete drifter.colonpipe.biz:bkup