# systemd/.config/systemd/user/rbackup.service
[Unit]
Description=Backs up files

[Service]
ExecStart=/usr/bin/unison-2.40 -auto -batch %h/./org ssh://drifter.colonpipe.biz/bkup/org