# kde/.config/autostart/xbindkeys.desktop
[Desktop Entry]
Exec=xbindkeys -p
Icon=
Name=xbindkeys
Path=
Terminal=False
Type=Application
