# kde/.config/autostart/xwaylandvideobridge.desktop
[Desktop Entry]
Exec=xwaylandvideobridge
Icon=
Name=xwaylandvideobridge
Path=
Terminal=False
Type=Application
